# roles/configure_zfs_datasets/templates/zfs-load-key.service.j2
[Unit]
Description=Load ZFS encryption keys and restore dependent services
DefaultDependencies=no
Before=zfs-mount.service
{% if zfs_restart_nfs_server | default(true) %}
Before=nfs-server.service
{% endif %}
After=zfs-import.target
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/zfs load-key -a
ExecStart=/sbin/zfs mount -a
{% if zfs_restore_iscsi_targets | default(true) %}
ExecStart=/bin/sleep 2
ExecStart=/usr/bin/targetctl restore
{% endif %}
{% if zfs_restart_nfs_server | default(true) %}
ExecStart=/sbin/zfs share -a
{% endif %}

[Install]
WantedBy=zfs-mount.service
